# octosip-api.service
[Unit]
Description=OctoSIP Monitor API
After=network.target postgresql.service

[Service]
EnvironmentFile=/opt/octosip/octosip.env
ExecStart=/opt/octosip/venv/bin/python3 ${OCTOSIP_DIR}/octosip_api.py
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
